# deployment/scripts/SRE~
#!/bin/sh

### BEGIN INIT INFO

### END INIT INFO

ROOT_PATH=/home/ziwei/workspace/KompicsSRE/deployment/
SRE_JAR=SRE-Kompics.jar
SRE_CONFIG=SREEnv.config

case "$1" in
  start)
	java -jar "$ROOT_PATH$SRE_JAR" "$ROOT_PATH" "$SRE_CONFIG" $
	;;
  stop)
	curl "http://localhost:8080/SRE/stop"
	;;
  restart|force-reload)
  	$0 stop
	sleep 8
	$0 start
	;;
  status)
	curl "http://localhost:8080/SRE/log"
	;;
  *)
	N=/etc/init.d/$NAME
	echo "Usage: $N {start|stop|restart|status}" >&2
	exit 1
	;;
esac

exit 0
